# crosstools-aarch64-gcc-9.2-linux-4.19-glibc-2.30-binutils-2.32/aarch64-buildroot-linux-gnu/sysroot/etc/profile.d/umask.sh
version https://git-lfs.github.com/spec/v1
oid sha256:9b7dae25ad0e172974b7d845a5d3d76e2f62a06b6556fd9c523031419c78d16a
size 10
